# root/etc/init.d/rc
#!/sbin/openrc-run

command=/usr/bin/caddy -agree -conf /config/Caddyfile -log stdout -root /tmp